# automation/com.user.transcript-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.transcript-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.14/bin/python3</string>
        <string>__PROJECT_DIR__/daily_sync.py</string>
        <string>--transcripts-dir</string>
        <string>__PROJECT_DIR__/transcripts</string>
        <string>--db-path</string>
        <string>__PROJECT_DIR__/processing.db</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>__PROJECT_DIR__/logs/sync.log</string>

    <key>StandardErrorPath</key>
    <string>__PROJECT_DIR__/logs/sync.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Library/Frameworks/Python.framework/Versions/3.14/bin</string>
        <key>ANTHROPIC_API_KEY</key>
        <string>__ANTHROPIC_API_KEY__</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
